# assets/services/burble.service
# SPDX-License-Identifier: MPL-2.0
#
# Burble Elixir/Phoenix control-plane service — *system* unit.
#
# Why a system unit (not --user) by default: Burble.Bolt.Listener binds
# udp/9 (WoL-compat poke port), which is a privileged port. systemd
# --user units silently ignore AmbientCapabilities= because user
# instances cannot grant capabilities the user doesn't already have;
# only PID 1's systemd can. A system unit with
# AmbientCapabilities=CAP_NET_BIND_SERVICE actually binds udp/9.
#
# Installed by scripts/install-service.sh (sudo) — @REPO_DIR@ and
# @USER@ are rewritten at install time.
#
# Logs: journalctl -u burble -f

[Unit]
Description=Burble voice/media control plane (Elixir, Bolt udp/7373+9, http :4020)
Documentation=https://github.com/hyperpolymath/burble
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=@USER@
Group=@USER@
WorkingDirectory=@REPO_DIR@/server
Environment=MIX_ENV=dev
Environment=LANG=C.UTF-8
ExecStart=/usr/bin/env mix phx.server
Restart=on-failure
RestartSec=5
TimeoutStopSec=15
KillMode=mixed
# This is the whole point of running as a system unit — grants the
# privileged-port bind to a non-root user.
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
NoNewPrivileges=true
StandardOutput=journal
StandardError=journal
SyslogIdentifier=burble

[Install]
WantedBy=multi-user.target
